# .docker-compose/supervisor/conf.d/websockets.conf
[program:websockets]
directory=/var/www/app/
process_name=%(program_name)s_%(process_num)02d
command=php artisan websockets:serve
numprocs=1
autostart=true
autorestart=true
user=root